# file/etc/uci-defaults/99-custom-settings.bak
#!/bin/sh

# 1. Password root & Izin File
echo "Setting root password..."
# Menggunakan cara yang lebih standar untuk ganti password di script
echo -e "root\nroot" | passwd root > /dev/null 2>&1

# Gunakan loop untuk chmod agar lebih rapi
for file in /sbin/free.sh /etc/rc.local; do
    [ -f "$file" ] && chmod +x "$file"
done

# 2. Konfigurasi UCI (System, Network, Firewall)
echo "Configuring UCI settings..."
uci batch <<EOF
    # System & Timezone
    set system.@system[0].hostname='FRDM-X'
    set system.@system[0].timezone='WIB-7'
    set system.@system[0].zonename='Asia/Jakarta'
    
    # NTP Servers
    delete system.ntp.server
    add_list system.ntp.server='0.id.pool.ntp.org'
    add_list system.ntp.server='1.id.pool.ntp.org'
    add_list system.ntp.server='time.google.com'
    
    # Network Interfaces
    set network.mm=interface
    set network.mm.proto='modemmanager'
    set network.mm.device='auto'
    set network.mm.apn='internet'
    set network.mm.peerdns='1'
    set network.mm.metric='10'
    
    set network.wan=interface
    set network.wan.proto='dhcp'
    set network.wan.device='eth1'
    
    set network.tethering=interface
    set network.tethering.proto='dhcp'
    set network.tethering.device='usb0'
    
    # Hapus wan6 jika ada
    delete network.wan6
    
    # Firewall (Pastikan interface masuk ke zona WAN/index 1)
    del_list firewall.@zone[1].network='wan6'
    add_list firewall.@zone[1].network='tethering'
    add_list firewall.@zone[1].network='wan'
    add_list firewall.@zone[1].network='mm'
    
    # Web UI (uhttpd & PHP)
    set uhttpd.main.index_page='index.php'
    set uhttpd.main.interpreter='.php=/usr/bin/php-cgi'
    
    commit network
    commit system
    commit firewall
    commit uhttpd
EOF

# 3. Optimasi & Fixes
echo "Applying optimizations..."

# Matikan signature check opkg (hanya jika belum dimatikan)
sed -i 's/^option check_signature/# option check_signature/g' /etc/opkg.conf 2>/dev/null

# USB Power Management
for dev in /sys/bus/usb/devices/*/power/control; do 
    [ -e "$dev" ] && echo "on" > "$dev"
done

# PHP Memory Limit
[ -f /etc/php.ini ] && sed -i 's/^memory_limit.*/memory_limit = 64M/' /etc/php.ini

# Ttyd Fix (Hanya jalan jika ttyd terinstall)
if [ -f /etc/config/ttyd ] && [ -f /bin/bash ]; then
    uci set ttyd.@ttyd[0].command='/bin/bash --login'
    uci commit ttyd
    /etc/init.d/ttyd enable
fi

# 4. Symlinks untuk TinyFM (Gunakan -f agar tidak error jika link sudah ada)
mkdir -p /www/tinyfm
ln -sf / /www/tinyfm/rootfs
ln -sf /etc/openclash/ /www/tinyfm/openclash
ln -sf /etc/openclash/ /www/tinyfm/nikki

# 5. Restart Services
echo "Restarting services..."
sysctl -p /etc/sysctl.conf >/dev/null 2>&1
/etc/init.d/network restart
/etc/init.d/uhttpd restart
/etc/init.d/ttyd restart
sleep 2
/etc/init.d/modemmanager restart

# Jalankan script tambahan jika ada
[ -f "/sbin/free.sh" ] && sh /sbin/free.sh
[ -f "/etc/custom_service/free.sh" ] && sh /etc/custom_service/free.sh

echo "Done! Router FRDM-X is ready."
exit 0
