# file/etc/uci-defaults/99-custom-settings
#!/bin/sh

# 1. Password root & Izin File
(echo "root"; sleep 1; echo "root") | passwd > /dev/null
[ -f "/usr/sbin/free.sh" ] && chmod 0755 /usr/sbin/free.sh
[ -f "/etc/rc.local" ] && chmod +x /etc/rc.local

# 2. Gabungan Semua Konfigurasi UCI
uci batch <<EOF
    # System & Timezone
    set system.@system[0].hostname='FRDM-X'
    set system.@system[0].timezone='WIB-7'
    set system.@system[0].zonename='Asia/Jakarta'
    delete system.ntp.server
    add_list system.ntp.server='0.id.pool.ntp.org'
    add_list system.ntp.server='time.google.com'
    
    # ModemManager WAN (Interface MM)
    delete network.wan
    delete network.wan6
    set network.mm=interface
    set network.mm.proto='modemmanager'
    set network.mm.device='auto'
    set network.mm.apn='internet'
    set network.mm.peerdns='1'
    set network.mm.metric='10'
    
    # Firewall, DHCP DNS, & TTYD
    add_list firewall.@zone[1].network='mm'
    add_list dhcp.lan.dns='8.8.8.8'
    add_list dhcp.lan.dns='1.1.1.1'
    set ttyd.@ttyd[0].enabled='1'
    set ttyd.@ttyd[0].command='/bin/login -f root'
    
    # ZRAM
    set zram.@zram[0].enabled='1'
    set zram.@zram[0].size='128'
    
    commit
EOF

# 3. Optimasi Eksternal (Non-UCI)
[ -f "/etc/opkg.conf" ] && sed -i 's/option check_signature/# option check_signature/g' /etc/opkg.conf
for dev in /sys/bus/usb/devices/usb*/power/control; do echo "on" > "$dev" 2>/dev/null; done
sysctl -p /etc/sysctl.conf

# 4. Jalankan Service & Script
/usr/sbin/free.sh &
/etc/init.d/modemmanager restart

exit 0
