# file/etc/uci-defaults/99-custom-settings
#!/bin/sh

# 1. Password root & Izin File
echo "Setting root password..."
(echo "root"; sleep 1; echo "root") | passwd > /dev/null
[ -f "/usr/sbin/free.sh" ] && chmod 0755 /usr/sbin/free.sh
[ -f "/etc/rc.local" ] && chmod +x /etc/rc.local

# 2. Konfigurasi UCI (Optimasi Network & Firewall)
uci batch <<EOF
    # System & Timezone
    set system.@system[0].hostname='FRDM-X'
    set system.@system[0].timezone='WIB-7'
    set system.@system[0].zonename='Asia/Jakarta'
    del system.ntp.server
    add_list system.ntp.server='0.id.pool.ntp.org'
    add_list system.ntp.server='1.id.pool.ntp.org'
    add_list system.ntp.server='2.id.pool.ntp.org'
    add_list system.ntp.server='3.id.pool.ntp.org'
    add_list system.ntp.server='time.google.com'
    
    # Network Interface
    set network.mm=interface
    set network.mm.proto='modemmanager'
    set network.mm.device='auto'
    set network.mm.apn='internet'
    set network.mm.peerdns='1'
    set network.mm.metric='10'
    
    set network.wan=interface
    set network.wan.proto='dhcp'
    set network.wan.device='eth1'
    
    set network.tethering=interface
    set network.tethering.proto='dhcp'
    set network.tethering.device='usb0'
    
    delete network.wan6
    
    # Firewall Zona (Menggunakan nama zona 'wan' agar lebih akurat)
    del_list firewall.wan.network='wan6'
    add_list firewall.wan.network='tethering'
    add_list firewall.wan.network='wan'
    add_list firewall.wan.network='mm'
    
    commit network
    commit system
    commit firewall
    commit dhcp
EOF

# 3. Optimasi Eksternal (Non-UCI)
sed -i 's/^option check_signature/# option check_signature/g' /etc/opkg.conf

# USB Power Management
for dev in /sys/bus/usb/devices/*/power/control; do 
    [ -e "$dev" ] && echo "on" > "$dev" 2>/dev/null
done

# Reload sysctl
sysctl -p /etc/sysctl.conf >/dev/null 2>&1

# 4. Jalankan Service Berat di Background
# Menggunakan '&' agar booting tidak tertahan proses pencarian sinyal modem
(
    /etc/init.d/network restart
    sleep 10
    /etc/init.d/modemmanager restart
    [ -f "/usr/sbin/free.sh" ] && /usr/sbin/free.sh
) &

# 5. Ttyd Fix (Agar Bash langsung jalan)
if [ -f /bin/bash ]; then
    uci set ttyd.@ttyd[0].command='/bin/bash --login'
    uci commit ttyd
    /etc/init.d/ttyd enable
    /etc/init.d/ttyd restart
fi

exit 0
